# assets/systemd/gunicorn-dev.socket
[Unit]
Description=gunicorn-dev socket

[Socket]
ListenStream=/run/gunicorn.sock

[Install]
WantedBy=sockets.target